# docker/supervisord.conf
[program:neo4j]
command=/usr/share/neo4j/bin/neo4j console
user=neo4j
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/neo4j.log
stderr_logfile=/var/log/supervisor/neo4j-error.log
priority=10

[program:strfry]
command=/usr/local/bin/strfry relay
user=strfry
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/strfry.log
stderr_logfile=/var/log/supervisor/strfry-error.log
environment=STRFRY_CONFIG="/etc/strfry.conf"
priority=20

[program:strfry-router]
command=/usr/local/bin/strfry router /etc/strfry-router-tapestry.config
user=strfry
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/strfry-router.log
stderr_logfile=/var/log/supervisor/strfry-router-error.log
environment=STRFRY_CONFIG="/etc/strfry.conf"
priority=25

[program:brainstorm]
command=/usr/local/bin/start-brainstorm.sh
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/brainstorm.log
stderr_logfile=/var/log/supervisor/brainstorm-error.log
environment=NODE_ENV="production",HOME="/root"
priority=30
startsecs=10
startretries=5
