# systemd/monitor_deprecated.service
[Unit]
Description=My Daily Script

[Service]
ExecStart=/usr/local/lib/node_modules/brainstorm/src/manage/monitor/monitor.sh